# templates/default/lifeguard-upstart.conf.erb
description "<%= @service.name || @service.service %>"

start on <%= @service.depends || "runlevel [2345]" %>
stop on runlevel [016]

<% if @service.dir %>
chdir <%= @service.dir %>
env HOME=<%= @service.dir %>
<% end %>

<% (@service.env||{}).each do |k,v| %>
env <%= k %>=<%= v %>
<% end %>

<% if @campfire %>
env CAMPFIRE_ACCOUNT=<%= @campfire[:account] %>
env CAMPFIRE_TOKEN=<%= @campfire[:token] %>
env CAMPFIRE_ROOM=<%= @campfire[:room] %>
<% end %>

console log

setuid <%= @service.user %>

script
<% if @service.path %>PATH=<%= @service.path %>:$PATH <% end %>/usr/bin/lifeguard <% if @service.monitor_dir %>--dir "<%= @service.monitor_dir %>" <% end %>--cmd "<%= @service.command %>" --title "<%= @service.name || @service.service %>" <% if @service.handoff %>--handoff<% end %>
end script

respawn
respawn limit <%= node.lifeguard.respawn_count %> <%= node.lifeguard.respawn_interval %>